# vulcanbar/systemd/vulcanbar.service
[Unit]
Description=VulcanBar Touch Bar Daemon
Documentation=https://github.com/VulcanOS/vulcanbar
After=graphical.target
ConditionPathExists=/dev/dri/card0

[Service]
Type=simple

# Legacy mode (tiny-dfr compatible) - uses /etc/vulcanbar/tiny-dfr.toml
ExecStart=/usr/bin/vulcanbar

# Modular mode (Waybar-style) - uses /etc/vulcanbar/config.toml
# Uncomment the line below and comment out the line above to use modular mode
#ExecStart=/usr/bin/vulcanbar --modular

Restart=on-failure
RestartSec=3

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
ReadWritePaths=/sys/class/backlight

# Required device access
SupplementaryGroups=input video

# Environment
# Set to debug for verbose logging: RUST_LOG=vulcanbar=debug
Environment=RUST_LOG=vulcanbar=info

[Install]
WantedBy=multi-user.target
