# dotfiles/systemd/.config/systemd/user/goimapnotify.service
[Unit]
Description=Go IMAP Notify (Gmail)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=goimapnotify -conf %h/.config/goimapnotify/goimapnotify.yaml -log-level info
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
